# scripts/homeai-gpu-recover.path
[Unit]
Description=Watch CDI spec for NVIDIA driver changes; trigger GPU container recovery
# Run AFTER nvidia-cdi-refresh has had a chance to write the new spec.
After=nvidia-cdi-refresh.service
Wants=nvidia-cdi-refresh.service

[Path]
# /var/run/cdi/nvidia.yaml is auto-regenerated by nvidia-cdi-refresh.path on
# every NVIDIA driver / nvidia-ctk install or upgrade.
PathChanged=/var/run/cdi/nvidia.yaml
Unit=homeai-gpu-recover.service

[Install]
WantedBy=multi-user.target
